# server/infra/faasta.service
[Unit]
Description=Faasta Function-as-a-Service Platform
After=network.target
Documentation=https://github.com/fourlexboehm/faasta

[Service]
Type=simple
User=faasta
Group=faasta
WorkingDirectory=/opt/faasta
ExecStart=/bin/sh -lc 'exec /opt/faasta/kvmserver --cwd /opt/faasta --allow-read=/lib --allow-read=/usr/lib --allow-read=/etc/ssl --allow-read=/etc/resolv.conf --allow-read=/etc/hosts --allow-read=/etc/nsswitch.conf --allow-read=/etc/gai.conf --allow-write=/opt/faasta --allow-listen --allow-connect storage --1-to-1 /opt/faasta/faasta-server storage-vm --db-path ./data/db --metrics-db-path ./data/metrics ++ run /opt/faasta/faasta-server serve --base-domain faasta.lol'
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=append:/var/log/faasta.log
StandardError=append:/var/log/faasta.error.log

# Limits
LimitNOFILE=65536
TimeoutStartSec=30
TimeoutStopSec=30

# Security options
ProtectSystem=full
PrivateTmp=true
ProtectHome=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
