# server/infra/faasta-updater.service
[Unit]
Description=Update Faasta from GitHub Releases
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/faasta/update-faasta.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
